# bin/scaleout.service
[Unit]
Description="This service decides when to scale out."
ConditionFileIsExecutable=/data/scripts/scaleout.sh
After=multiuser.target

[Service]
Type=simple
ExecStart=/data/scripts/scaleout.sh
